# etc/service/dnsmasq/run
#!/bin/sh

exec /usr/sbin/dnsmasq -d -x /var/run/dnsmasq/dnsmasq.pid -u dnsmasq -r /var/run/dnsmasq/resolv.conf -7 /etc/dnsmasq.d,.dpkg-dist,.dpkg-old,.dpkg-new >>/var/log/dnsmasq.log 2>&1
